# files/internals/sim.cron
MAILTO=
SHELL=/bin/sh
*/4 * * * * root /usr/local/sbin/sim -q >> /dev/null 2>&1
